# debian/kerrighed-tools.postinst
#!/bin/sh

set -e

update-rc.d kerrighed defaults > /dev/null

##DEBHELPER##
